# scripts/llm-wiki-hub.timer
[Unit]
Description=Run LLM-Wiki-Hub aggregate every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
